# root/etc/services.d/irssi/run
#!/usr/bin/with-contenv bash

_term() {
  echo "Caught SIGTERM signal!"
  killall -TERM irssi 2>/dev/null
}

trap _term SIGTERM

screen -D -m -S irssi \
  s6-setuidgid abc \
  /usr/bin/irssi --home /config/.irssi &

wait
